# motd/update-motd.d/25-reboot-required
#!/bin/bash

if [ -f /var/run/reboot-required ]; then
    echo -e "\033[1;31m*** System restart required ***\033[0m"
fi
